# deploy/cometblue-mcp.service.j2
[Unit]
Description=CometBlue Control — MCP HTTP/SSE server
After=network.target bluetooth.target cometblue.service
Wants=bluetooth.target

[Service]
Type=simple
User={{ service_user }}
WorkingDirectory={{ install_dir }}
ExecStart={{ install_dir }}/.venv/bin/cometblue-control mcp --transport http --host 0.0.0.0 --port {{ mcp_port }}
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cometblue-mcp

[Install]
WantedBy=multi-user.target
