# deploy/cometblue.service.j2
[Unit]
Description=CometBlue Control — BLE thermostat management service
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User={{ service_user }}
WorkingDirectory={{ install_dir }}
ExecStart={{ install_dir }}/.venv/bin/cometblue-control serve
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cometblue

[Install]
WantedBy=multi-user.target
